# fg06/fg06.service
[Unit]
Description="FG06"
[Service]
ExecStart=/usr/bin/init_usb
Type=simple
[Install]
WantedBy=multi-user.target
